# systemd/acp-broker.service
# Authoritative copy of the systemd --user unit for the per-devvm
# acp-broker. Symlinked into ~/.config/systemd/user/ by ~/dotfiles/init.sh
# (Linux block); after editing this file, pick up changes with:
#
#   systemctl --user daemon-reload
#   systemctl --user restart acp-broker
#   ss -lxp | grep acp-broker.sock        # expect: acp-broker
#
# Companion to ~/dotfiles/launchd/com.mkarrmann.persistence-server.plist
# on the Mac side: the broker connects out to that persistence-server
# over the reverse tunnel established by nvs-tunnels (-r 7847:7847).
#
# Unit specifiers used:
#   %h  $HOME
#   %t  $XDG_RUNTIME_DIR  (matches the UDS path acp-broker-up used)
#   %S  $XDG_STATE_HOME   (~/.local/state by default)
#   %l  short hostname    (used as ACP_BROKER_ID for per-client metadata)

[Unit]
Description=ACP broker (one per devvm; brokers ACP traffic from per-session nvims)
After=network.target

[Service]
Type=simple
Environment=ACP_BROKER_ID=%l
Environment=ACP_BROKER_SQLITE_PERSISTENCE_ADDR=127.0.0.1:7847
Environment=DVSC_AUTO_RUN_MODE=autoRunAll
# Verbose tracing across every broker-internal crate (broker core, plugin
# SDK, sqlite-persistence-plugin/persistence-plugin, persistence-server
# client) so resume/load issues can be diagnosed from the log alone
# without rebuilding. Notable INFO-level events become visible: dispatch
# ladder branches in `client_link::dispatch_load_session`, `try_resume`
# return values, `replay_session_history_to_client`, agent capability
# lookups, capture/upload pipeline. `info` is kept for everything else
# to avoid spamming hyper/tokio/tonic internals.
# To temporarily raise to trace for a single crate at runtime:
#   systemctl --user set-environment RUST_LOG=acp_broker=trace,...
#   systemctl --user restart acp-broker
Environment=RUST_LOG=info,acp_broker=debug,acp_broker_plugin_sdk=debug,sqlite_persistence_plugin=debug,persistence_plugin=debug,persistence_server_protocol=debug,ralph_loop_plugin=debug
# Backtraces in panic logs — invaluable when the broker crashes mid-RPC
# and the only surviving artifact is server.log.
Environment=RUST_BACKTRACE=1
StateDirectory=acp-broker
# Defensive: remove a stale UDS file left behind by a previous broker that
# died without running Drop (SIGKILL/OOM/hard kill). Without this, every
# restart fails with EADDRINUSE and `Restart=always` loops forever. Safe
# because systemd serializes Stop -> ExecStartPre -> ExecStart, so the
# previous broker is fully terminated by the time this runs. Remove once
# acp-broker's `UdsListener::bind` retries after a connect-probe.
ExecStartPre=/bin/rm -f %t/acp-broker.sock
ExecStart=%h/.cargo/bin/acp-broker --declare-agent name=claude,cmd=%h/bin/claude-agent-acp,default=true,native_fork=claude --declare-agent name=dvsc-core,cmd=%h/bin/dvsc-core-acp,native_fork=claude --socket %t/acp-broker.sock --per-session-agents --allow-agent-spawn
Restart=always
RestartSec=5
StandardOutput=append:%S/acp-broker/server.log
StandardError=append:%S/acp-broker/server.log

[Install]
WantedBy=default.target
