# systemd/omnigent-hub-reconcile.service
[Unit]
Description=Reconcile Omnigent hub services with shared ownership
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PATH=%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
EnvironmentFile=%h/dotfiles/omnigent_config/topology.env
ExecStart=%h/dotfiles/bin/omnigent-hub-reconcile
TimeoutStartSec=180
StandardOutput=append:%S/omnigent-hub/reconcile.log
StandardError=append:%S/omnigent-hub/reconcile.log
UMask=0077
