# systemd/omnigent-snapshot.timer
[Unit]
Description=Snapshot active Omnigent state every five minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
RandomizedDelaySec=15
Unit=omnigent-snapshot.service

[Install]
WantedBy=timers.target
