# systemd/omnigent-host.service
# Authoritative systemd --user unit for the Omnigent execution host.
# Symlinked into ~/.config/systemd/user/ by ~/dotfiles/init.sh (Linux block).
#
#   systemctl --user restart omnigent-host
#   omnigent host status --server "$(omnigent-server-url)" --json
#
# Registers THIS devserver as an execution host on the Omnigent server (so
# agents' tool calls run where the code lives). The server URL is resolved per
# host by ~/dotfiles/bin/omnigent-server-url: loopback on the HUB, the HUB's
# FQDN over the prod network on other devservers. Runs on every devserver; the
# HUB also runs omnigent-server.service.

[Unit]
Description=Omnigent execution host (registers this devserver with the hub server)
After=network.target omnigent-server.service
Wants=omnigent-server.service

[Service]
Type=simple
Environment=PATH=%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
EnvironmentFile=%h/dotfiles/omnigent_config/topology.env
# Agents spawned by the host need env that systemd --user does NOT inherit from
# the login shell:
#  - x509 for mTLS to Meta's AI gateway (used by codex; claude uses api-key-helper).
#  - CODEX_MANAGED_CONFIG_PATH=/dev/null: the system codex config at
#    /etc/codex/managed_config.toml sets model_provider=responses without
#    defining the provider, so codex fails at startup unless it's skipped
#    (~/.localrc does this for interactive shells). No x2p proxy is needed here:
#    unlike the Mac, this devserver is on-corp and reaches the gateway directly.
Environment=THRIFT_TLS_CL_CERT_PATH=/var/facebook/credentials/%u/x509/%u.pem
Environment=THRIFT_TLS_CL_KEY_PATH=/var/facebook/credentials/%u/x509/%u.pem
Environment=CODEX_MANAGED_CONFIG_PATH=/dev/null
Environment=XDG_RUNTIME_DIR=%t
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
# The host filters its environment before launching runners. Forward only the
# non-secret paths/selectors required by Meta's Claude and Codex launchers.
Environment=OMNIGENT_RUNNER_ENV_PASSTHROUGH=THRIFT_TLS_CL_CERT_PATH,THRIFT_TLS_CL_KEY_PATH,CODEX_MANAGED_CONFIG_PATH,XDG_RUNTIME_DIR,DBUS_SESSION_BUS_ADDRESS
ExecStart=%h/bin/omnigent-host-run
Restart=always
RestartSec=5
StateDirectory=omnigent-host
StandardOutput=append:%S/omnigent-host/host.log
StandardError=append:%S/omnigent-host/host.log

[Install]
WantedBy=default.target
