# systemd/omnigent-prodnet.service
# Authoritative systemd --user unit for the Omnigent prod-network proxy.
# Symlinked into ~/.config/systemd/user/ by ~/dotfiles/init.sh (Linux block).
#
# The Omnigent server binds 127.0.0.1 only (a v6 bind is IPv6-only under uvicorn
# and would break the v4 loopback the Mac ET forwards + Orchest use). This proxy
# exposes the hub's prod-network v6 address :6767 and relays to the loopback
# server, so the OTHER devserver can dial the hub's FQDN directly. HUB only.
#
#   systemctl --user restart omnigent-prodnet
#   curl --noproxy '*' http://$(hostname -f):6767/health   # expect 200

[Unit]
Description=Omnigent prod-network proxy (v6 :6767 -> loopback server; hub only)
After=network.target omnigent-server.service
Wants=omnigent-server.service

[Service]
Type=simple
Environment=PATH=%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
EnvironmentFile=%h/dotfiles/omnigent_config/topology.env
ExecCondition=%h/dotfiles/services/omnigent-hub/.venv/bin/omnigent-hub gate
ExecStart=/usr/bin/env python3 %h/bin/omnigent-prodnet-proxy
Restart=always
RestartSec=5
TimeoutStartSec=150
StateDirectory=omnigent-prodnet
StandardOutput=append:%S/omnigent-prodnet/proxy.log
StandardError=append:%S/omnigent-prodnet/proxy.log

[Install]
WantedBy=default.target
