# launchd/com.mkarrmann.omnigent-host.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Authoritative launchd job for the Mac-local Omnigent host daemon.

  omnigent-host-ensure starts the foreground host inside a persistent tmux or
  screen session because a plain nohup-launched host exits on this machine.
  The helper is idempotent and StartInterval lets launchd repair the host if it
  exits after login.

  Activation is automatic via ~/dotfiles/init.sh.

  Quick verify:
    launchctl print gui/$UID/com.mkarrmann.omnigent-host | head
    omnigent host status --server http://127.0.0.1:6767 --json
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.mkarrmann.omnigent-host</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/mkarrmann/bin/omnigent-host-ensure</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>60</integer>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/Users/mkarrmann/bin:/Users/mkarrmann/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>OMNIGENT_PORT</key>
      <string>6767</string>
      <!-- Codex reaches Meta's AI gateway through the local x2p proxy over mTLS.
           launchd does not inherit the interactive shell's exports (~/.zshrc), so
           a host-spawned codex subprocess needs these explicitly or it can't
           route (HTTP 421 "no upstream configured for this host"). Claude
           (claude-sdk) is unaffected — it auths via api-key-helper. -->
      <key>X2P_AGENT_PROXY_ADDRESS</key>
      <string>localhost:10054</string>
      <key>THRIFT_TLS_CL_CERT_PATH</key>
      <string>/var/facebook/credentials/mkarrmann/x509/mkarrmann.pem</string>
      <key>THRIFT_TLS_CL_KEY_PATH</key>
      <string>/var/facebook/credentials/mkarrmann/x509/mkarrmann.pem</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mkarrmann/.local/state/omnigent-host/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mkarrmann/.local/state/omnigent-host/stderr.log</string>
  </dict>
</plist>
